# deploy/systemd/edge-health.service
[Unit]
Description=reverse_ssh VPS edge health reporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/reverse-logger/edge-health.env
ExecStart=/usr/local/bin/edge-health
Restart=always
RestartSec=5
User=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/run /sys/class/net
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
