# deploy/systemd/edge-logger.service
[Unit]
Description=reverse_ssh VPS edge TCP proxy logger
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/reverse-logger/edge-logger.env
ExecStart=/usr/local/bin/edge-logger
Restart=always
RestartSec=5
DynamicUser=true
StateDirectory=reverse-logger
StateDirectoryMode=0750
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
